# gadget/usb-midi-ecm-gadget.service
[Unit]
Description=DEV USB gadget: MIDI + USB-Ethernet (SSH-over-USB) via configfs
After=sysinit.target
Before=fsr-midi.service
ConditionPathExists=/sys/class/udc

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/usb_midi_ecm_gadget.sh

[Install]
WantedBy=multi-user.target
