# systemd/fsr-midi.service
[Unit]
Description=FSR pressure sensor -> USB MIDI
# Start as early as possible; MIDI is USB-only so we depend on NOTHING network-related.
DefaultDependencies=no
After=systemd-modules-load.service
# Only run once the ADC/IIO device is present (BB-ADC overlay).
ConditionPathExists=/sys/bus/iio/devices/iio:device0
# Bring us up as part of normal boot but without pulling in extra ordering.
Before=multi-user.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/fsr-midi/fsr_midi.py
Restart=on-failure
RestartSec=1
Nice=-10
# Optional real-time scheduling for tighter trigger timing (uncomment if needed):
# CPUSchedulingPolicy=fifo
# CPUSchedulingPriority=50
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
